# eee.conf
[program:eee]
command=python3 /home/curator2/eee/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator2/eee/test.err.log
stdout_logfile=/home/curator2/eee/test.out.log